# frontend/cfh-frontend.service
[Unit]
Description=CFH Project - Frontend React Application
After=network.target cfh-project.service

[Service]
Type=simple
User=laraveluat
Group=laraveluat
WorkingDirectory=/opt/lampp/htdocs/MIMS/centralized-flow-hub/frontend
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PORT=3000"
ExecStart=/usr/bin/serve -s build -p 3000

# Restart policy
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/cfh-project/frontend-access.log
StandardError=append:/var/log/cfh-project/frontend-error.log

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
